# remote/remote-shared/deploy/remote-compose.service
[Unit]
Description=Remote auth/admin Linux deployment stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/mydw
ExecStart=/usr/bin/docker compose -f /opt/mydw/remote/remote-shared/deploy/docker-compose.linux.yml up -d --build
ExecStop=/usr/bin/docker compose -f /opt/mydw/remote/remote-shared/deploy/docker-compose.linux.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
